# ops/systemd/priceai-api-transit-public.service
[Unit]
Description=PriceAI API transit public pricing collector
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/priceai-nonshop
ExecStart=/opt/priceai-nonshop/run-api-transit-public.sh
TimeoutStartSec=240
