# ops/shop-collectors/systemd/priceai-hot-offer-verifier-heyuan.timer
[Unit]
Description=Run PriceAI Heyuan hot offer verifier every 5 minutes

[Timer]
OnCalendar=*-*-* *:03,08,13,18,23,28,33,38,43,48,53,58:00
AccuracySec=10s
RandomizedDelaySec=10
Persistent=false
Unit=priceai-hot-offer-verifier.service

[Install]
WantedBy=timers.target
